# nanosuspend.service
[Unit]
Description=Suspend script
[Service]
Type=simple
ExecStart=/usr/local/bin/nanosuspend.sh
[Install]
WantedBy=multi-user.target
